# deploy/live-runner.service
[Unit]
Description=Gemini BTC binary-option market maker (live)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/gtc-trading
EnvironmentFile=/etc/gtc-trading/live.env
ExecStart=/home/ubuntu/gtc-trading/.venv/bin/python /home/ubuntu/gtc-trading/src/strategy/runner.py --config /home/ubuntu/gtc-trading/config/live.json
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/ubuntu/gtc-trading/logs/live_runner.stdout.log
StandardError=append:/home/ubuntu/gtc-trading/logs/live_runner.stderr.log

# Hardening
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
