# deploy/paper-runner.service
[Unit]
Description=GTC Gemini BTC paper trader
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/gtc-2
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ubuntu/gtc-2/.venv/bin/python -u src/strategy/paper_runner.py --config config/paper.ec2.json
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
